# systemd/system/abyss-backup-sessions-checksum.timer
[Unit]
Description=Write daily checksum manifest for Codex and AoA session evidence when ABYSSVAULT is mounted
Documentation=file:{{ABYSS_BACKUP_ROOT}}/AGENTS.md

[Timer]
OnCalendar=*-*-* 22:10:00
RandomizedDelaySec=30m
Persistent=false
Unit=abyss-backup-sessions-checksum.service

[Install]
WantedBy=timers.target
