# systemd/system/abyss-observability-collect.timer
[Unit]
Description=Run Abyss Machine observability collector every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=60s
AccuracySec=15s
Persistent=true
Unit=abyss-observability-collect.service

[Install]
WantedBy=timers.target
