# systemd/system/abyss-ai-workload-refresh.service
[Unit]
Description=Abyss AI workload safe refresh
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md

[Service]
Type=oneshot
ExecStart={{ABYSS_LOCAL_BIN_DIR}}/abyss-machine ai workload refresh --json
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
