# templates/systemd/system/abyss-backup-system.timer
[Unit]
Description=Run ABYSSVAULT system backup twice weekly when the vault is mounted
Documentation=file:{{ABYSS_BACKUP_ROOT}}/AGENTS.md

[Timer]
OnCalendar=Wed,Sun *-*-* 03:30:00
RandomizedDelaySec=30m
Persistent=false
Unit=abyss-backup-system.service

[Install]
WantedBy=timers.target
