# templates/systemd/user/abyss-machine-typing-atspi-text-events.service
[Unit]
Description=Abyss Machine live committed text-change intake
Documentation=file://{{ABYSS_MACHINE_STATE}}/typing/AGENTS.md
After=graphical-session.target

[Service]
Type=simple
ExecStart={{ABYSS_LOCAL_BIN_DIR}}/abyss-machine typing atspi-text-events --forever --json
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target
