# templates/systemd/user/abyss-nervous-semantic-maintain.service
[Unit]
Description=Abyss Nervous semantic embedding maintenance
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/nervous/AGENTS.md
After=abyss-nervous-index-build.service

[Service]
Type=oneshot
Nice=13
IOSchedulingClass=idle
CPUWeight=25
IOWeight=20
MemoryHigh=768M
TimeoutStartSec=2100s
Environment=ABYSS_NERVOUS_SEMANTIC_MAINTAIN_TRIGGER=timer
ExecStart={{ABYSS_LOCAL_BIN_DIR}}/abyss-machine nervous semantic-maintain --json --no-thermal-sample
StandardOutput=null
StandardError=journal
