# templates/systemd/user/abyss-nervous-browser-content-capture.timer
[Unit]
Description=Run Abyss Nervous Firefox page-content capture frequently while Firefox is open
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/nervous/AGENTS.md

[Timer]
OnBootSec=2min
OnStartupSec=2min
OnUnitActiveSec=120s
AccuracySec=30s
Persistent=false

[Install]
WantedBy=timers.target
