# templates/systemd/user/abyss-machine-typing-saved-text-scan.timer
[Unit]
Description=Abyss Machine saved typed-text scan timer
Documentation=file://{{ABYSS_MACHINE_STATE}}/typing/AGENTS.md

[Timer]
OnBootSec=120s
OnUnitActiveSec=120s
AccuracySec=30s
Persistent=false

[Install]
WantedBy=timers.target
