# templates/systemd/user/abyss-nervous-derived-refresh.service
[Unit]
Description=Abyss Nervous derived refresh after passive chronicle snapshot
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/nervous/AGENTS.md
After=abyss-nervous-passive-chronicle.service

[Service]
Type=oneshot
Nice=12
IOSchedulingClass=idle
CPUWeight=30
IOWeight=20
MemoryHigh=512M
TimeoutStartSec=360s
Environment=ABYSS_NERVOUS_DERIVED_TRIGGER=passive-chronicle
ExecStart={{ABYSS_LOCAL_BIN_DIR}}/abyss-machine resource launch --class medium --kind indexing --unattended --timeout 300 --success-on-block --no-thermal-sample --json -- {{ABYSS_LOCAL_BIN_DIR}}/abyss-machine nervous quality-audit --refresh --json
ExecStart={{ABYSS_LOCAL_BIN_DIR}}/abyss-machine resource launch --class medium --kind indexing --unattended --timeout 300 --success-on-block --no-thermal-sample --json -- {{ABYSS_LOCAL_BIN_DIR}}/abyss-machine nervous index-build --json
StandardOutput=null
StandardError=journal
