# templates/systemd/user/abyss-gemma4-spark.timer
[Unit]
Description=Ensure Abyss Gemma4 Spark resident candidate is started when gates allow
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/ai/llm/AGENTS.md

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=abyss-gemma4-spark.service

[Install]
WantedBy=timers.target
