# templates/systemd/user/abyss-mirror.service
[Unit]
Description=Abyss repo mirror (sanitized)

[Service]
Type=oneshot
ExecStart=/srv/abyss-git/scripts/mirror.sh
