# templates/systemd/user/abyss-machine-doctor.timer
[Unit]
Description=Run Abyss Machine doctor self-maintenance
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/doctor/AGENTS.md

[Timer]
OnBootSec=25min
OnUnitActiveSec=2h
AccuracySec=15min
Persistent=true

[Install]
WantedBy=timers.target
