# templates/systemd/user/abyss-machine-typing-codex-session-tail.service
[Unit]
Description=Abyss Machine Codex submitted prompt tail intake
Documentation=file://{{ABYSS_MACHINE_STATE}}/typing/AGENTS.md

[Service]
Type=simple
ExecStart={{ABYSS_LOCAL_BIN_DIR}}/abyss-machine typing codex-session-tail --files 8 --interval 1 --forever
Restart=always
RestartSec=2s

[Install]
WantedBy=default.target
