# systemd/user/abyss-nervous-index-build.service
[Unit]
Description=Abyss nervous-system local index rebuild
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/nervous/AGENTS.md

[Service]
Type=oneshot
Nice=10
IOSchedulingClass=idle
CPUWeight=30
IOWeight=20
MemoryHigh=512M
TimeoutStartSec=360s
Environment=ABYSS_NERVOUS_INDEX_TRIGGER=timer
ExecStart={{ABYSS_LOCAL_BIN_DIR}}/abyss-machine resource launch --class medium --kind indexing --unattended --timeout 300 --success-on-block --json -- {{ABYSS_LOCAL_BIN_DIR}}/abyss-machine nervous index-build --json
StandardOutput=null
StandardError=journal
