# systemd/user/abyss-process-snapshot.service
[Unit]
Description=Abyss Machine lightweight process snapshot
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/processes/AGENTS.md

[Service]
Type=oneshot
Nice=10
IOSchedulingClass=idle
ExecStart={{ABYSS_LOCAL_BIN_DIR}}/abyss-machine processes snapshot --top 60 --json
StandardOutput=null
StandardError=journal
