# systemd/user/abyss-desktop-session-fd-guard.service
[Unit]
Description=Abyss desktop session fd/pidfd guard
Documentation=file:{{ABYSS_MACHINE_STATE}}/processes/AGENTS.md
After=graphical-session.target

[Service]
Type=oneshot
ExecStart={{ABYSS_USER_HOME}}/.local/bin/abyss-desktop-session-fd-guard --apply --quiet
Nice=10
IOSchedulingClass=idle
