# systemd/user/abyss-nervous-passive-chronicle.timer
[Unit]
Description=Run Abyss Nervous passive chronicle snapshots periodically
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=2min
Persistent=true

[Install]
WantedBy=timers.target
