# systemd/user/abyss-machine-typing-codex-session-tail.timer
[Unit]
Description=Abyss Machine Codex submitted prompt tail intake timer
Documentation=file://{{ABYSS_MACHINE_STATE}}/typing/AGENTS.md

[Timer]
OnBootSec=15s
OnUnitActiveSec=5s
AccuracySec=1s
Persistent=false

[Install]
WantedBy=timers.target
