# systemd/user/abyss-dictation-server.service
[Unit]
Description=Abyss Dictation warm OpenVINO server
After=default.target

[Service]
Type=simple
ExecStart={{ABYSS_LOCAL_LIBEXEC_DIR}}/abyss-dictation-server --socket %t/abyss-machine/dictation/server.sock --preload-model-dir {{ABYSS_OS_ROOT}}/abyss-stack/Models/stt/whisper/openvino/whisper-large-v3-turbo --preload-device AUTO:GPU,CPU --preload-language ru --warmup
Restart=on-failure
RestartSec=2s
MemoryLow=768M
MemoryHigh=4096M
Environment=OPENVINO_LOG_LEVEL=2
Environment=ABYSS_OPENVINO_CACHE_DIR={{ABYSS_MACHINE_SRV}}/cache/ai/openvino
Environment=HF_HOME={{ABYSS_MACHINE_SRV}}/cache/ai/huggingface
Environment=XDG_CACHE_HOME={{ABYSS_MACHINE_SRV}}/cache/home/%u/cache
Environment=TRANSFORMERS_OFFLINE=1
Environment=HF_HUB_OFFLINE=1
Environment=TOKENIZERS_PARALLELISM=false

[Install]
WantedBy=default.target
