# systemd/user/abyss-machine-typing-nervous-refresh.service
[Unit]
Description=Abyss Machine typing nervous refresh tick
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/typing/AGENTS.md
After=abyss-machine-typing-codex-session-tail.service abyss-machine-typing-saved-text-scan.service

[Service]
Type=oneshot
Nice=12
IOSchedulingClass=idle
CPUWeight=25
IOWeight=20
MemoryHigh=768M
TimeoutStartSec=420s
Environment=ABYSS_TYPING_NERVOUS_REFRESH_TRIGGER=timer
ExecStart={{ABYSS_USER_HOME}}/.local/bin/abyss-machine-typing-nervous-refresh-tick
StandardOutput=null
StandardError=journal
