# systemd/user/abyss-gemma4-spark-monitor.timer
[Unit]
Description=Record Abyss Gemma4 Spark resident status every minute
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/ai/llm/AGENTS.md

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=abyss-gemma4-spark-monitor.service

[Install]
WantedBy=timers.target
