# systemd/user/abyss-nervous-semantic-maintain.timer
[Unit]
Description=Run Abyss Nervous semantic embedding maintenance
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/nervous/AGENTS.md

[Timer]
OnBootSec=20min
OnUnitActiveSec=90min
AccuracySec=10min
Persistent=true

[Install]
WantedBy=timers.target
