# systemd/user/abyss-gemma4-spark-jobs.timer
[Unit]
Description=Run Abyss Gemma4 Spark resident useful jobs periodically
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/ai/llm/AGENTS.md

[Timer]
OnBootSec=20min
OnUnitActiveSec=2h
AccuracySec=5min
RandomizedDelaySec=20min
Persistent=true
Unit=abyss-gemma4-spark-jobs.service

[Install]
WantedBy=timers.target
