# systemd/user/abyss-gemma4-spark-digest.service
[Unit]
Description=Abyss Gemma4 Spark dictation digest tick
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/ai/llm/AGENTS.md

[Service]
Type=oneshot
Nice=10
IOSchedulingClass=idle
CPUWeight=20
IOWeight=20
MemoryHigh=256M
TimeoutStartSec=90s
Environment=ABYSS_GEMMA4_SPARK_SERVING_OWNER=abyss-stack
Environment=ABYSS_GEMMA4_SPARK_BASE_URL=http://127.0.0.1:11435
Environment=ABYSS_GEMMA4_SPARK_ALLOW_LOCAL_FALLBACK=0
ExecStart={{ABYSS_MACHINE_SRV}}/tools/abyss-gemma4-spark-resident digest --limit 3 --json
StandardOutput=null
StandardError=journal
