# systemd/user/abyss-gemma4-spark.service
[Unit]
Description=Abyss Gemma4 Spark stack endpoint bridge
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/ai/llm/AGENTS.md
After=default.target
StartLimitIntervalSec=10min
StartLimitBurst=4

[Service]
Type=simple
Nice=8
IOSchedulingClass=best-effort
IOSchedulingPriority=6
CPUWeight=40
IOWeight=40
MemoryLow=256M
MemoryHigh=6144M
Environment=XDG_CACHE_HOME={{ABYSS_MACHINE_SRV}}/cache/home/%u/cache
Environment=HF_HOME={{ABYSS_MACHINE_SRV}}/cache/ai/huggingface
Environment=TOKENIZERS_PARALLELISM=false
Environment=OMP_NUM_THREADS=4
Environment=MKL_NUM_THREADS=4
Environment=OPENBLAS_NUM_THREADS=4
Environment=NUMEXPR_NUM_THREADS=4
Environment=ABYSS_GEMMA4_SPARK_SERVING_OWNER=abyss-stack
Environment=ABYSS_GEMMA4_SPARK_BASE_URL=http://127.0.0.1:11435
Environment=ABYSS_GEMMA4_SPARK_ALLOW_LOCAL_FALLBACK=0
ExecStart={{ABYSS_MACHINE_SRV}}/tools/abyss-gemma4-spark-resident serve
Restart=on-failure
RestartSec=30s
TimeoutStopSec=45s

[Install]
WantedBy=default.target
