# systemd/user/abyss-nervous-browser-content-capture.service
[Unit]
Description=Abyss Nervous Firefox page-content capture tick
Documentation=file://{{ABYSS_MACHINE_ETC}}/AGENTS.md
Documentation=file://{{ABYSS_MACHINE_STATE}}/nervous/AGENTS.md

[Service]
Type=oneshot
Nice=10
IOSchedulingClass=idle
CPUWeight=10
IOWeight=10
MemoryHigh=256M
TimeoutStartSec=8s
ExecStart=%h/.local/bin/abyss-nervous-browser-content-capture-tick
StandardOutput=null
StandardError=journal
